# deploy/systemd/aethos.service
[Unit]
Description=Aethos ACP messaging bridge
Documentation=https://github.com/aesoteric/aethos
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=aethos
Group=aethos
Environment=HOME=/var/lib/aethos/home
Environment=AETHOS_DATA_DIR=/var/lib/aethos/data
EnvironmentFile=-/etc/aethos/aethos.env
WorkingDirectory=/var/lib/aethos/workspace
ExecStart=/usr/local/bin/aethos
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/aethos
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
